# deploy/vak.service
[Unit]
Description=vāk — Voice-first thinking partner
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=vak
Group=vak
WorkingDirectory=/home/vak/app/backend
ExecStart=/home/vak/app/backend/.venv/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --ws-max-size 5242880
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vak

[Install]
WantedBy=multi-user.target
